# init/profile.d/rust.sh
#!/usr/bin/env bash
#shellcheck disable=1090

[[ -f "${HOME}/.cargo/env" ]] && source ${HOME}/.cargo/env